# templates/default/nodejs.upstart.conf.erb
# <%= @instance %> Upstart script
# Generated by Chef for <%= @fqdn %>

description "Pushit NodeJS Application Server"

start on runlevel [2345]
stop on runlevel [06]

respawn

<% @env.each do |key, value| -%>
env <%= key %>="<%= value %>"
<% end -%>

setuid <%= @user %>
setgid <%= @group %>

chdir <%= @app_path %>
exec start-stop-daemon \
  --start \
  --chdir <%= @app_path %> \
  --chuid <%= @user %>:<%= @group %> \
  --make-pidfile --pidfile <%= @pid_path %> \
  --exec <%= @exec %> \
  -- \
  <%= @script_path %> \
  >> <%= @log_path %> 2>&1
