# templates/default/nodejs.upstart.conf.erb
# <%= @instance %> Upstart script
# Generated by Chef for <%= @fqdn %>

description "Pushit NodeJS Application"

export PATH="<%= @env_path %>:$PATH"

start on (local-filesystems and runlevel [2345])
stop on runlevel [06]

respawn
console log

setuid root
setgid root

chdir <%= @app_path %>

script
  export HOME="/opt/pushit"
<% @env.each do |key, value| -%>
  export <%= key %>="<%= value %>"
  <% end -%>

  echo $$ > <%= @pid_file %>
  exec <%= @exec %> <%= @script_path %> >> <%= @log_path %> 2>&1
end script

pre-start script
  echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> <%= @log_path %>
end script

pre-stop script
  rm <%= @pid_file %>
  echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> <%= @log_path %>
end script