# templates/default/rails.upstart.conf.erb
# <%= @instance %> Upstart script
# Generated by Chef for <%= @fqdn %>

description "Pushit Rails Application Server"

export PATH="<%= @env_path %>:$PATH"

start on (local-filesystems and runlevel [2345])
stop on runlevel [06]

respawn

<% @env.each do |key, value| -%>
env <%= key %>="<%= value %>"
<% end -%>

exec start-stop-daemon \
  --start \
  --chuid <%= @user %>:<%= @group %> \
  --make-pidfile --pidfile <%= @pid_file %> \
  --exec <%= @exec %> \
  -- \
  -E <%= @environment %> \
  -c <%= @config_file %> \
  >> <%= @log_file %> 2>&1
